# docker/crontab
0 3 * * * /usr/local/bin/abrechnung db clean
